# dev-libs/pd-mapper/files/pd-mapper.init
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Qualcomm PD mapper service"
pidfile="/run/${SVCNAME}.pid"
command="/usr/bin/pd-mapper"
command_background=true

depend() {
    need qrtr-ns
    use logger
}
